# wifiap/etc/udev/rules.d/90-wlan0-ap.rules
KERNEL=="wlan0", ACTION=="add", RUN+="/usr/bin/bash /usr/bin/setup-wlan0.sh"
